# examples/example/output/services/dtmi_mqtt_broker_command_service_1.service

[Unit]
Description=Your Custom Service
[Service]
ExecStart=mosquitto -c /conf/mosquitto.conf
Restart=always
User=root
Group=root
[Install]
WantedBy=multi-user.target